# icecream-server.service
[Unit]
Description=Ice Cream Paradise Server
After=network.target mongodb.service

[Service]
Type=simple
User=jaivik
WorkingDirectory=/home/jaivik/Documents/icecream-project
ExecStart=/usr/bin/node server/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=append:/home/jaivik/Documents/icecream-project/logs/server.log
StandardError=append:/home/jaivik/Documents/icecream-project/logs/error.log

[Install]
WantedBy=multi-user.target

